# Library/LaunchAgents/com.drewdawson.claude.phase-d-cost-control-loop.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.drewdawson.claude.phase-d-cost-control-loop</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/drewdawson/.claude/scripts/phase_d_cost_control_loop.sh</string>
  </array>
  <key>StartInterval</key><integer>900</integer>
  <key>RunAtLoad</key><true/>
  <key>StandardOutPath</key><string>/Users/drewdawson/.claude/reports/phase-d-cost-control-loop.launchd.out.log</string>
  <key>StandardErrorPath</key><string>/Users/drewdawson/.claude/reports/phase-d-cost-control-loop.launchd.err.log</string>
</dict>
</plist>
